# etc/network/if-up.d/udhcp
#!/bin/sh
[ "$IFACE" != "lo" ] || exit 0
udhcpc -t 5 -q -s /etc/dhcp.script